# scripts/systemd/verse-nightjar.timer.in
[Unit]
Description=Daily Verse Nightjar morning edition

[Timer]
OnCalendar=*-*-* 04:30:00 UTC
AccuracySec=5m
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
